# images/k3s-node/ansible/files/bcachefs-mount.service
[Unit]
Description = Mount Storage disks

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStop=/usr/bin/bcachefs-mount.sh

[Install]
WantedBy=multi-user.target
